# simplefan.service
[Unit]
Description=Fan Control
After=syslog.target

[Service]
ExecStart=/usr/bin/simplefan

[Install]
WantedBy=multi-user.target
